# deploy/ec2/veritas-pocketbase.service
[Unit]
Description=Veritas PocketBase Document Store
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/veritas-ph
ExecStart=/home/ubuntu/veritas-ph/pb_bin/pocketbase serve --http 127.0.0.1:8090
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
