# deploy/ec2/veritas-crawler.service
[Unit]
Description=Veritas Background Ingestion Crawler (PhilGEPS & PDF Scraping)
After=network.target veritas-api.service

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/veritas-ph/apps/api
ExecStart=/home/ubuntu/veritas-ph/.venv_linux/bin/python workers/crawler_worker.py
Restart=always
RestartSec=10
EnvironmentFile=/home/ubuntu/veritas-ph/apps/api/.env

[Install]
WantedBy=multi-user.target
