# deploy/ec2/veritas-api.service
[Unit]
Description=Veritas FastAPI Server
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/veritas-ph/apps/api
ExecStart=/home/ubuntu/veritas-ph/.venv_linux/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
EnvironmentFile=/home/ubuntu/veritas-ph/apps/api/.env

[Install]
WantedBy=multi-user.target
