# deploy/ec2/veritas-worker.service
[Unit]
Description=Veritas Background Worker (Crawlers & Analyzers)
After=network.target veritas-api.service

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/veritas-ph/apps/api
ExecStart=/home/ubuntu/veritas-ph/.venv_linux/bin/python workers/local_worker.py
Restart=always
RestartSec=10
EnvironmentFile=/home/ubuntu/veritas-ph/apps/api/.env

[Install]
WantedBy=multi-user.target
